# backups/cron_files/crontab
0 0 * * * root /bin/pbdb_backup.sh >> /var/log/pbdb_backup.log 2>&1
